# dot_config/systemd/user/udiskie.service
[Unit]
Description=Udiskie daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/udiskie --smart-tray
Restart=always
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target
