# dot_config/systemd/user/waypaper.service
[Unit]
Description=Set a random wallpaper with waypaper
Documentation=man:waypaper(1)

[Service]
Type=oneshot
ExecStart=waypaper --random
