# dot_config/systemd/user/waybar.service.d/override.conf
[Service]
Type=exec
# Repurpose XDG Autostart filtering
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "Hyprland" ""
Slice=app-graphical.slice
